# scripts/keep-warm/keep-warm.service
[Unit]
Description=Recommend-a-book API Keep-Warm Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nobody
# Replace with your actual path and API URL
ExecStart=/path/to/recommend-a-book/scripts/keep-warm/keep-warm.sh https://recommend-a-book-api.onrender.com
Restart=always
RestartSec=60
Environment="PING_INTERVAL=600"
Environment="LOG_FILE=/var/log/keep-warm.log"

# Optional: uncomment to enable verbose logging
# Environment="VERBOSE=true"

[Install]
WantedBy=multi-user.target
